# taiyi.service
[Unit]
Description=Taiyi AGI Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=nicola
Group=nicola
WorkingDirectory=/home/nicola/.openclaw/workspace

# Clash 代理环境变量
Environment="http_proxy=http://127.0.0.1:7890"
Environment="https_proxy=http://127.0.0.1:7890"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# 启动命令
ExecStart=/usr/bin/python3 -c "from skills.zhiji.monitor import print_dashboard; print_dashboard()"
Restart=always
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=taiyi

[Install]
WantedBy=multi-user.target
